# tttt/conf/cron/greentea.cron
# Green Tea
# https://github.com/SatelliteQE/GreenTea/wiki/Setup
# Copy file to directory /etc/cron.d/ and restart service cron

# For asynchronous operation Green Tea needs to run cron
*/1 * * * * 	greentea 	python /data/Greantea/manage.py pickup --traceback

# Following command check status of beaker jobs (automation tests)
*/20 * * * * 	greentea 	python /data/Greantea/manage.py check --quiet --traceback
